# launchagents/com.srouter.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- __SROUTER_LAUNCHAGENT_MARKER__ -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__SROUTER_LAUNCHAGENT_LABEL__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__SROUTER_PYTHON_BIN__</string>
    <string>__SROUTER_DASHBOARD_PATH__</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__SROUTER_ROOT_DIR__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
